# rootfs/common/usr/lib/systemd/system/systemd-growfs@sysroot.service.d/override.conf
[Service]
MountFlags=slave
ExecStartPre=mount -o remount,rw /sysroot