# rootfs/common/etc/systemd/system/tailscaled.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/sbin/tailscaled --state=/var/lib/tailscale/tailscaled.state --socket=/run/tailscale/tailscaled.sock
